# 20170520/104 redis
wget http://download.redis.io/releases/redis-3.2.9.tar.gz
tar xzf redis-3.2.9.tar.gz

mv redis-3.2.9 /opt

cd /opt/redis-3.2.9
#yum -y install gcc
make
#如果报找不到jmalloc
#make MALLOC=libc



The binaries that are now compiled are available in the src directory. Run Redis with:

src/redis-server



启动过程有报错和警告，按照提示进行处理
vi /etc/sysctl.conf 
sysctl vm.overcommit_memory=1

echo '' >> /etc/rc.locl
echo 1000 >/proc/sys/net/core/somaxconn



You can interact with Redis using the built-in client:
src/redis-cli
redis> set foo bar
OK
redis> get foo
"bar"






注意点：
ctrl+z 会把服务挂起来，导致客户端无法访问，应将服务启动方式改为后台启动
redis.conf文件将daemonize改成yes

redis-server & 启动






开机启动配置


1. 根据启动脚本要求，将修改好的配置文件以端口为名复制一份到指定目录。需使用root用户。
mkdir /etc/redis
cp redis.conf /etc/redis/6379.conf
2. 将启动脚本复制到/etc/init.d目录下，本例将启动脚本命名为redisd（通常都以d结尾表示是后台自启动服务）。
cp redis_init_script /etc/init.d/redisd
3.  设置为开机自启动
此处直接配置开启自启动 chkconfig redisd on 将报错误： service redisd does not support chkconfig
参照 此篇文章 ，在启动脚本开头添加如下两行注释以修改其运行级别：
#!/bin/sh
# chkconfig:   2345 90 10
# description:  Redis is a persistent key-value database
#
再设置即可成功。
#设置为开机自启动服务器
chkconfig redisd on
#打开服务
service redisd start
#关闭服务
service redisd stop













